# systemd/watchclaw.timer
[Unit]
Description=Run WatchClaw periodic host watchdog scan
Documentation=https://github.com/jcestevezr/watchclaw

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=watchclaw.service

[Install]
WantedBy=timers.target
