# systemd/watchclaw.service
[Unit]
Description=WatchClaw host watchdog scan
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=@WATCHCLAW_BIN@ run-once --config @WATCHCLAW_CONFIG@
