# systemd/watchclaw.service
[Unit]
Description=WatchClaw watchdog service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/env watchclaw status
Restart=on-failure

[Install]
WantedBy=multi-user.target
