# systemd/watchclaw.service
[Unit]
Description=WatchClaw host watchdog scan
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/env python3 -m watchclaw.cli run-once --config @WATCHCLAW_CONFIG@
